# src/ezrassor_rover/Systemctl Starter files/ros2-start-forwarder.service
[Unit]
description=Starts ROS2 nodes automatically

[Service]
User=root
Group=root
ExecStart=/etc/systemd/user/ros2_start_forwarder.sh

[Install]
WantedBy=multi-user.target

